# amneziawg-tools-pkg/wg-quick.target
[Unit]
Description=WireGuard Tunnels via awg-quick(8)
